# httpd_cookbook/templates/apache-exporter.service.erb
[Unit]
Description=Prometheus Apache Exporter
Documentation=https://github.com/Lusitaniae/apache_exporter
After=network.target

[Service]
User=prometheus
ExecStart=/usr/local/bin/apache_exporter \
          --telemetry.path=<%= @telemetry_path %> \
          --scrape_uri=<%= @scrape_uri %>
Restart=on-failure
RestartSec=5s

# Security enhancements
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target